# roles/rsync-backup/templates/sslbackup-full.timer.j2
[Unit]
Description=Run rsyncd-clean daily at {{ sslbackup_time }}

[Timer]
OnCalendar=Sun *-*-* {{ sslbackup_time }}

[Install]
WantedBy=timers.target

