# roles/rsync-backup/templates/sslbackup-diff.service.j2
[Unit]
Description=rsyncd diff backup

[Service]
Type=oneshot
ExecStart= /usr/bin/sslbackup.sh diff
