# roles/rsync-backup/templates/sslbackup-full.service.j2
[Unit]
Description=rsyncd full backup.

[Service]
Type=oneshot
ExecStart= /usr/bin/sslbackup.sh full
