# old_roles/rsyncd/templates/rsyncd-clean.service.j2
[Unit]
Description=Clean rsyncd.

[Service]
Type=oneshot
ExecStart= /usr/bin/rsyncd-clean.sh
